# bot/vds-agent.service
[Unit]
Description=VDS Secure AI Agent
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/vds-agent
Restart=always
User=root
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/socks-monitor/vds-agent.env

[Install]
WantedBy=multi-user.target
